# docker/home_index.ini
[program:django]
command=gunicorn home_index.wsgi:application -b :8000 --reload
directory=/src/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /src/log/supervisor/django.log
loglevel=info
stopasgroup=true
killasgroup=true
