# docker/nginx.ini
[program:nginx]
command=nginx -g "daemon off;"
directory=/src
user=root
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile = /data/log/supervisor/git.log
loglevel=info
stopasgroup=true
killasgroup=true